# roles/kde-desktop-customisation/templates/configure_plasma.service
[Unit]
After=xdg-desktop-autostart.target

[Install]
WantedBy=xdg-desktop-autostart.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'echo >&4 "executed on $(date)" '


