# roles/kde-desktop-customisation/templates/configure_plasma.service
[Unit]
After=pulseaudio.service

[Install]
WantedBy=default.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'echo >&4 "executed on $(date)" qdbus org.kde.plasmashell /PlasmaShell org.kde.PlasmaShell.evaluateScript "$(</home/{{ desktop_username }}/.config/ansible_distro/configure_plasma_shell.js)"'