# roles/kde-desktop-customisation/templates/configure_plasma.service
[Unit]
After=xdg-desktop-autostart.target

[Install]
WantedBy=xdg-desktop-autostart.target

[Service]
Type=simple
ExecStart=/home/{{ desktop_username }}/.config/ansible_distro/configure_plasma_shell.sh


